# deploy/systemd/stage.api.service
[Unit]
Description=GarPr Stage API service

[Service]
ExecStart=/usr/local/bin/twistd \
    --nodaemon \
    --pidfile= \
    --logfile="logs/api.log" \
    -oy serve_api.tac

WorkingDirectory=/home/deploy/stage/garpr

User=root
Group=root

Restart=always

[Install]
WantedBy=multi-user.target
